# crontab
# Hacker Experience Legacy - Crontab (PHP 8 unified stack)
# Update GAME_PATH to match your deployment directory
GAME_PATH=/home/openclaw/he/legacy-master
PHP=/usr/local/bin/php
LOG=/var/log/game/cron.log

# === Every minute ===
* * * * *			$PHP	$GAME_PATH/cron/doomUpdater.php		>> $LOG 2>&1
* * * * *			$PHP	$GAME_PATH/cron/newRoundUpdater.php	>> $LOG 2>&1

# === Every 5 minutes ===
*/5 * * * *			$PHP	$GAME_PATH/cron/fbiUpdate.php		>> $LOG 2>&1

# === Every 10 minutes ===
05,15,25,35,45,55 * * * *	$PHP	$GAME_PATH/cron/defcon2.php		>> $LOG 2>&1
05,15,25,35,45,55 * * * *	$PHP	$GAME_PATH/cron/endWar2.php		>> $LOG 2>&1

# === Every 15 minutes ===
05,20,35,50 * * * *		$PHP	$GAME_PATH/cron/removeExpiredLogins.php	>> $LOG 2>&1
*/15 * * * *			$PHP	$GAME_PATH/cron/restoreNPC.php		>> $LOG 2>&1
*/15 * * * *			$PHP	$GAME_PATH/cron/generateMissions.php	>> $LOG 2>&1

# === Every 30 minutes ===
*/30 * * * *			$PHP	$GAME_PATH/cron/safenetUpdate.php	>> $LOG 2>&1
19,49 * * * *			$PHP	$GAME_PATH/cron/badgeHunter.php		>> $LOG 2>&1

# === Hourly ===
00 * * * *			$PHP	$GAME_PATH/cron/updateCurStats.php	>> $LOG 2>&1
00 * * * *			$PHP	$GAME_PATH/cron/updateRanking.php	>> $LOG 2>&1
00 * * * *			$PHP	$GAME_PATH/cron/rankGenerator.php	>> $LOG 2>&1
02 * * * *			$PHP	$GAME_PATH/cron/removeDownNPC.php	>> $LOG 2>&1
09 * * * *			$PHP	$GAME_PATH/cron/removeExpiredAccs.php	>> $LOG 2>&1
09 * * * *			$PHP	$GAME_PATH/cron/removeExpiredPremium.php	>> $LOG 2>&1
10 * * * *			$PHP	$GAME_PATH/cron/updateServerStats.php	>> $LOG 2>&1
12 * * * *			$PHP	$GAME_PATH/cron/queryCounter.php		>> $LOG 2>&1
19 * * * *			$PHP	$GAME_PATH/cron/removeExpiredNPC.php	>> $LOG 2>&1
19 * * * *			$PHP	$GAME_PATH/cron/badgeHunterAll.php	>> $LOG 2>&1
29 * * * *			$PHP	$GAME_PATH/cron/backup_game.php		>> $LOG 2>&1
40 * * * *			$PHP	$GAME_PATH/cron/removeExpiredHTMLPages.php >> $LOG 2>&1

# === Every 3 hours ===
10 */3 * * *			$PHP	$GAME_PATH/cron/antivirusNPC.php	>> $LOG 2>&1
